# trunk/target/linux/x86/alix2/base-files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 76
/svn/hiwifi/!svn/ver/3/trunk/target/linux/x86/alix2/base-files/etc/hotplug.d
END
