# trunk/package_bak/broadcom-wl/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 68
/svn/hiwifi/!svn/ver/3/trunk/package/broadcom-wl/files/etc/hotplug.d
END
