# trunk/package_bak/qos-scripts/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 68
/svn/hiwifi/!svn/ver/3/trunk/package/qos-scripts/files/etc/hotplug.d
END
