# trunk/target/linux/atheros/base-files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 74
/svn/hiwifi/!svn/ver/3/trunk/target/linux/atheros/base-files/etc/hotplug.d
END
