# trunk/target/linux/omap24xx/base-files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 75
/svn/hiwifi/!svn/ver/3/trunk/target/linux/omap24xx/base-files/etc/hotplug.d
END
